# init/libvirt-bridges.conf
# libvirt-bridges - configure network bridhes for libvirt-bin
#

description	"configure bridges for virtual network devices"


start on starting libvirt-bin

task

exec /usr/sbin/virsh-bridges start
